# deploy/booking-backend.service
[Unit]
Description=Booking Backend API
After=redis.service

[Service]
Type=simple
WorkingDirectory=/home/backup/upgrade
Environment=PYTHONPATH=/home/backup/upgrade
Environment=PATH=/home/backup/upgrade/venv/bin:/usr/bin
EnvironmentFile=/home/backup/upgrade/.env
ExecStart=/home/backup/upgrade/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
